# ansible/playbooks/roles/webserver/files/systemd/acme-renew-certificate.timer
[Unit]
Description=Renew certificates with acme.sh

[Timer]
OnCalendar=*-*-* 19:07:12
Persistent=true

[Install]
WantedBy=timers.target
