# ansible/playbooks/roles/deploykit_monitor/files/conf/systemd/deploykit-reload-nginx.path
[Unit] 
Description=Reload Nginx when DeployKit receives new deployment

[Path] 
PathChanged=/srv/web/deploykit
Unit=deploykit-reload-nginx.service

[Install]
WantedBy=multi-user.target
